# packaging/deb/prerm
#!/bin/sh
set -e

case "$1" in
    remove|upgrade|deconfigure)
        if command -v systemctl >/dev/null 2>&1; then
            systemctl stop cam-shim.service 2>/dev/null || true
            systemctl disable cam-shim.service 2>/dev/null || true
        fi
        ;;

    failed-upgrade)
        ;;

esac

#DEBHELPER#

exit 0
